# templates/nimbus-eth1-short-benchmark.service.j2
[Unit]
Description=Nimbus Eth1 short import benchmarking service (24hrs)
Documentation=https://github.com/status-im/nimbus-eth1
Requires=network-online.target
After=network-online.target

[Service]
User={{ nimbus_eth1_service_user }}
Group={{ nimbus_eth1_service_group }}
SyslogIdentifier=nimbus-eth1-{{ nimbus_eth1_network }}-short-benchmark
StandardOutput=file:{{ nimbus_eth1_benchmark_log_output_file_path }}
StandardError=file:{{ nimbus_eth1_benchmark_log_error_file_path }}
SyslogLevel=debug
RuntimeMaxSec=24hr
TimeoutStartSec=36000
ExecStartPre={{ nimbus_eth1_build_script }}
ExecStartPre={{ nimbus_eth1_db_gen_script }}
ExecStartPre={{ nimbus_eth1_benchmark_gen_script }}
ExecStart={{ nimbus_eth1_repo_path }}/build/nimbus_execution_client import \
    --network={{ nimbus_eth1_network }} \
    --max-blocks=1000000 \
    --data-dir='{{ nimbus_eth1_data_path }}/{{ nimbus_eth1_data_folder }}' \
    --era1-dir={{ nimbus_eth1_era1_dir | mandatory }} \
    --era-dir={{ nimbus_eth1_era_dir | mandatory }} \
    --debug-csv-stats={{nimbus_eth1_benchmark_git_repo_dir}}/{{ nimbus_eth1_benchmark_file_name }}
ExecStopPost={{ nimbus_eth1_publish_script }}

[Install]
WantedBy=multi-user.target
