# templates/nimbus-eth1-short-benchmark.timer.j2
[Unit]
Description=Timer for Nimbus Eth1 short import benchmarking
After={{ nimbus_eth1_service_name }}.service

[Timer]
OnUnitInactiveSec=45min
RemainAfterElapse=true

[Install]
WantedBy=multi-user.target
