# templates/nimbus-eth1-benchmark.service.j2
[Unit]
Description=Nimbus Eth1 {{ nimbus_eth1_benchmark_type }} import benchmarking service ({{ "24hrs" if nimbus_eth1_benchmark_type == "short" else "1 week" }})
Documentation=https://github.com/status-im/nimbus-eth1
Requires=network-online.target
After=network-online.target

[Service]
User={{ nimbus_eth1_service_user }}
Group={{ nimbus_eth1_service_group }}
SyslogIdentifier=nimbus-eth1-{{ nimbus_eth1_network }}-{{ nimbus_eth1_benchmark_type }}-benchmark
SyslogFacility=local6
RuntimeMaxSec={{ "24hr" if nimbus_eth1_benchmark_type == "short" else "168hr" }}
TimeoutStartSec=36000
Environment=BENCHMARKING_TYPE={{ nimbus_eth1_benchmark_type }}
Environment=FORCE_RUN="false"
ExecStart="{{ nimbus_eth1_service_scripts_base }}/run-benchmark.sh"
Restart=always
RestartSec=60

[Install]
WantedBy=multi-user.target
