# templates/nimbus-eth1-long-benchmark.timer.j2
[Unit]
Description=Timer for Nimbus ETH1 Long Benchmark Service
Requires={{ nimbus_eth1_service_name }}.service

[Timer]
OnCalendar=*:0/30
AccuracySec=1s
Persistent=true

[Install]
WantedBy=timers.target
