# templates/nimbus-eth1-short-benchmark.timer.j2
[Unit]
Description=Timer for Nimbus ETH1 Short Benchmark Service
Requires=nimbus-eth1-short-benchmark.service

[Timer]
OnCalendar=*:0/30
AccuracySec=1s
Persistent=true

[Install]
WantedBy=timers.target
